# setup/colonies.cron
5 5 * * *   root    /usr/sbin/reboot
* * * * *   root    bash /home/pi/colonies/setup/check-tunnel.sh > /dev/null 2>&1

